# crontab
## OpenClaw Market Scanner - every 10 minutes
## flock -n prevents concurrent runs; source /etc/environment for proxy + API key
*/3 * * * * root flock -n /tmp/scanner.lock bash -c '. /etc/environment; cd /app && /usr/local/bin/node market_scanner.mjs >> /var/log/scanner.log 2>&1'
# Daily Report (6:00 AM Beijing time => 22:00 UTC)
0 22 * * * root bash -c '. /etc/environment; /usr/local/bin/node /app/daily_report.mjs >> /var/log/report.log 2>&1'
